# deb/src/main/resources/control/postinst
#!/bin/sh

# Ensure our service is enabled
mv /etc/bfm/minipg/minipg-app*.jar /etc/bfm/minipg/minipg-app.jar

systemctl enable minipg.service >/dev/null 2>&1

